# cron/app.crontab
#
MAILTO = ""
SHELL = /bin/bash
PATH = $HOME/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * ${HOME}/bin/chkrun.sh ${HOME}/bin/testapp
